# etc/systemd/system/daynightlocation.timer
[Unit]
Description=Start OAP DayNightLocation service 5 min after boot

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
